# dist/grzegorz_webui.service
[Unit]
Description=Grzegorz webUI

#Requires=grzegorz@.service

[Service]
ExecStart=grzegorz-webui

[Install]
WantedBy=default.target
